# personal_codex/automations/daily-skill-friction/com.hoteng.codex.daily-skill-friction-prefetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hoteng.codex.daily-skill-friction-prefetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/hoteng/Program/GitHub/Joey-Tools/codex-workspace/scripts/codex_workspace.py</string>
    <string>prefetch</string>
    <string>--stamp</string>
    <string>daily-skill-friction</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/hoteng/Program/GitHub/Joey-Tools/codex-workspace</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>50</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/hoteng/.codex/personal-sync/logs/daily-skill-friction-prefetch.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hoteng/.codex/personal-sync/logs/daily-skill-friction-prefetch.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
